# s6-overlay/s6-rc.d/laravel-queue-rabbitmq/run
#!/command/execlineb -P
with-contenv
cd /var/www/html
php
artisan
rabbitmq:consume-user-events