# s6-overlay/s6-rc.d/laravel-scheduler/run
#!/command/execlineb -P
with-contenv
cd /var/www/html
loopwhilex
foreground { php artisan schedule:run }
s6-sleep 60