# s6-overlay/s6-rc.d/laravel-queue-events/run
#!/command/execlineb -P
with-contenv
cd /var/www/html
php
artisan
queue:work
--queue=user_events_processing
--tries=3
--backoff=3
--max-jobs=1000
--max-time=3600
--timeout=60